# ops/deploy/phase-f1/systemd/pm2-root.service.d/phase-f1-handoff.conf
[Service]
# A single checksummed, atomically replaced state envelope decides boot ownership.
# PM2 is eligible only before the durable `freezing` transition.
ExecCondition=/usr/bin/node @@PHASE_F1_PACK_DIR@@/boot-owner-condition.mjs pm2
